# recipe/root/files/etc/udev/rules.d/99-kbd.rules
ACTION=="add", ATTRS{manufacturer}=="Iceman", RUN+="/usr/bin/kbd_exec.sh"
